# build/etc/services.d/tmagent/run
#!/usr/bin/with-contenv sh

set -e

# define the data and log dirs

DATADIR=${DATADIR:-data}
LOGDIR=${LOGDIR:-log}

export TMAGENT_DATADIR="/tmdata/$DATADIR" TMAGENT_LOGDIR="/tmdata/$LOGDIR"

exec /etc/ssstm/tmagent/tmagent
