# build/etc/services.d/tmlicd/run
#!/usr/bin/with-contenv sh

set -e

/usr/local/bin/tmlicd
